# deploy/systemd/sigap-web.service
[Unit]
Description=Sigap Next.js Application
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/apps/sigap-planner/web
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_URL=https://sigap.jasamaritim.co.id/api
Environment=PORT=3000
Environment=NEXT_TELEMETRY_DISABLED=1
ExecStart=/bin/npm start -- --hostname 0.0.0.0
Restart=on-failure
RestartSec=5
LimitNOFILE=10000

[Install]
WantedBy=multi-user.target
